# http.service
[Unit]
Description=A Minimal HTTP server (netcat, inline command)
Documentation=man:nc(1)
After=network.target
Wants=network-online.target

[Service]
Type=exec
DynamicUser=true
Environment=PORT=8080

# Pre-flight checks: fail early if nc binary or bash are missing.
ExecStartPre=/usr/bin/test -x /bin/bash
ExecStartPre=/bin/bash -lc 'command -v nc >/dev/null'

# ExecStart: inline bash loop that serves HTTP responses.
ExecStart=/bin/bash -c '\
  while true; do \
    echo -ne "HTTP/1.1 200 OK\r\nContent-Length: 12\r\nContent-Type: text/plain; charset=utf-8\r\nConnection: close\r\n\r\nHello World!" \
      | nc -l "$PORT"; \
  done'

# https://www.freedesktop.org/software/systemd/man/latest/systemd.service.html#Restart=
Restart=on-failure

[Install]
Alias=http.service
WantedBy=multi-user.target
